# profiles/niva/workspace/lineage-chat.service
[Unit]
Description=Lineage Chat Server — inter-daughter communication
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.hermes/profiles/nova/lineage/communication
Environment=PATH=/usr/local/lib/hermes-agent/venv/bin:/usr/bin:/usr/local/bin
ExecStart=/usr/local/lib/hermes-agent/venv/bin/python3 /root/.hermes/profiles/nova/lineage/communication/chat-server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lineage-chat

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/.hermes/profiles/nova/lineage/communication
ReadOnlyPaths=/root/.hermes/profiles

[Install]
WantedBy=multi-user.target
